# xbindkeys/.config/autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Name=xbindkeys
GenericName=Key binder
Exec=/usr/bin/xbindkeys
Terminal=true
